# luci-app-iptv-manager/root/etc/init.d/iptv-manager
#!/bin/sh /etc/rc.common
# ==========================================
# IPTV Manager init script — procd format
# Compatible with OpenWRT 25.x
# ==========================================

START=99
STOP=10
USE_PROCD=1

IPTV_MANAGER="/usr/share/iptv-manager/IPTV-Manager.sh"
IPTV_DIR="/etc/iptv"
IPTV_PORT="8082"
WWW_DIR="/www/iptv"
CGI_DIR="/www/iptv/cgi-bin"
HTTPD_PID="/var/run/iptv-httpd.pid"

# Read port from config if available
[ -f "$IPTV_DIR/defaults.conf" ] && . "$IPTV_DIR/defaults.conf"
IPTV_PORT="${IPTV_PORT:-8082}"

start_service() {
    # Check if playlist exists
    [ -f "$IPTV_DIR/playlist.m3u" ] || return 1

    # Prepare directories
    mkdir -p "$WWW_DIR" "$CGI_DIR" /www/cgi-bin

    # Copy playlist
    cp "$IPTV_DIR/playlist.m3u" "$WWW_DIR/playlist.m3u" 2>/dev/null || true

    # Generate CGI scripts from modular script
    if [ -f "$IPTV_MANAGER" ]; then
        . "$IPTV_MANAGER" --server 2>/dev/null || true
    fi

    # Kill any existing instance
    kill "$(cat "$HTTPD_PID" 2>/dev/null)" 2>/dev/null || true
    kill "$(pgrep -f "uhttpd.*$IPTV_PORT")" 2>/dev/null || true
    sleep 1

    # Start uhttpd via procd
    procd_open_instance iptv-manager
    procd_set_param command uhttpd \
        -f \
        -p "0.0.0.0:$IPTV_PORT" \
        -h "$WWW_DIR" \
        -x "$CGI_DIR" \
        -i ".cgi=/bin/sh" \
        -R

    procd_set_param pidfile "$HTTPD_PID"
    procd_set_param respawn 3600 5 0
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param env IPTV_MANAGER_MODE=daemon

    procd_close_instance

    # Record start time
    date +%s > /tmp/iptv-started 2>/dev/null || true
}

stop_service() {
    kill "$(cat "$HTTPD_PID" 2>/dev/null)" 2>/dev/null || true
    kill "$(pgrep -f "uhttpd.*$IPTV_PORT")" 2>/dev/null || true
    rm -f "$HTTPD_PID"
    rm -f /tmp/iptv-started
}

restart() {
    stop
    sleep 1
    start
}

reload() {
    if [ -f "$IPTV_MANAGER" ]; then
        . "$IPTV_MANAGER" --server 2>/dev/null || true
    fi
    stop
    sleep 1
    start
}

status() {
    if [ -f "$HTTPD_PID" ] && kill -0 "$(cat "$HTTPD_PID" 2>/dev/null)" 2>/dev/null; then
        echo "running"
    elif pgrep -f "uhttpd.*$IPTV_PORT" >/dev/null 2>&1; then
        echo "running"
    else
        echo "stopped"
    fi
}
